# infrastructure/systemd/ai-avangard-backend.service
[Unit]
Description=AI-Avangard Backend API Service
After=network.target postgresql.service redis.service qdrant.service

[Service]
Type=simple
User=temrjan
Group=temrjan
WorkingDirectory=/home/temrjan/ai-avangard
Environment="PATH=/home/temrjan/ai-avangard/venv/bin"
ExecStart=/home/temrjan/ai-avangard/venv/bin/uvicorn backend.app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --log-level info

Restart=always
RestartSec=5

StandardOutput=append:/var/log/ai-avangard/backend.log
StandardError=append:/var/log/ai-avangard/backend-error.log

[Install]
WantedBy=multi-user.target
